# deploy/init/supervisor-breeze-core.conf
# Breeze Core — supervisord program (any OS, including containers).
# Install:  cp deploy/init/supervisor-breeze-core.conf /etc/supervisor/conf.d/breeze-core.conf
#           supervisorctl reread && supervisorctl update
#
# No systemd sandbox — see HARDENING.md §7 "Hardening without systemd".

[program:breeze-core]
# Bind 127.0.0.1 instead of the LAN IP if a reverse proxy fronts it.
command=/opt/breeze-core/venv/bin/uvicorn meow_ac.app:app --host 192.168.1.10 --port 8420
directory=/opt/breeze-core
user=breeze
environment=AC_CONFIG="/etc/breeze-core/config.json",AC_DEVICES="/etc/breeze-core/devices.json",AC_PROGRAMS="/etc/breeze-core/programs.json"
autostart=true
autorestart=true
startsecs=5
; uvicorn shuts down cleanly on SIGINT
stopsignal=INT
stdout_logfile=/var/log/breeze-core.log
redirect_stderr=true
